# scripts/update.timer
[Unit]
Description = Timer to update System

[Timer]
OnBootSec=3min
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target






